# ansible/roles/restic/templates/restic-backup.timer.j2
[Unit]
Description=Restic Backup Timer

[Timer]
OnCalendar={{ restic_backup_schedule }}
Persistent=true
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
